# autocleaner.timer
[Unit]
Description=autocleaner timer

[Timer]
Unit=autocleaner.service
OnCalendar=*-*-* 23:55:00

[Install]
WantedBy=timers.target

