# autocleaner.service
[Unit]
Description=autocleaner

[Service]
Type=oneshot
ExecStart=/opt/autocleaner/main.rb

